# init.d/v2ray
#!/sbin/openrc-run

name="V2Ray"
description="A platform for building proxies to bypass network restrictions"
description_checkconfig="Test configuration file"

: ${env:="V2RAY_LOCATION_ASSET=/usr/local/lib/v2ray/"}
: ${conf:="/usr/local/etc/v2ray/config.json"}

command="/usr/local/bin/v2ray"
command_args="run -config $conf"
command_test_args="test -config $conf"
command_user="nobody"

pidfile="/run/v2ray.pid"
command_background="yes"

extra_commands="checkconfig"

depend() {
	need net
}

checkconfig() {
	export $env
	$command $command_test_args
}

start_pre() {
	checkconfig
}
